# automation/daily-sync/com.vamsios.daily-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.vamsios.daily-sync</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/vamsi/Coding/vamsios/automation/daily-sync/run.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>19</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/vamsi/Coding/vamsios/automation/daily-sync/launchd-stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/vamsi/Coding/vamsios/automation/daily-sync/launchd-stderr.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/vamsi/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>HOME</key>
		<string>/Users/vamsi</string>
	</dict>

	<key>WorkingDirectory</key>
	<string>/Users/vamsi/Coding/vamsios</string>

	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
