# launchd/jobs/com.vamsi.daily-brain-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vamsi.daily-brain-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/vamsi/Coding/vamsios/automation/daily-sync/run.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/vamsi/Coding/vamsios/automation/daily-sync/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vamsi/Coding/vamsios/automation/daily-sync/sync.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
