# launchd/jobs/com.motion.bootcamp-moderator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.motion.bootcamp-moderator</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/vamsi/Coding/Community Management/sweep.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SLACK_BOOTCAMP_TOKEN</key>
        <string>REPLACE_ME</string>
        <key>ANTHROPIC_API_KEY</key>
        <string>REPLACE_ME</string>
        <key>SWEEP_MINUTES</key>
        <string>60</string>
    </dict>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/Users/vamsi/Coding/Community Management/sweep.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vamsi/Coding/Community Management/sweep.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
